# tools/tools/etc/profile.d/ccm-tools-tomcat.cmd.in
if "%CCM_ZIP_ROOT%" == "" (
    set CCM_ZIP_ROOT=C:\ccm
)

set prefix=%CCM_ZIP_ROOT%\usr

set CCM_SERVLET_TOMCAT_CONF=%prefix%\share\ccm-tools\server\tomcat\bin\tomcat-conf-gen

set prefix=

