# tools/tools/etc/profile.d/ccm-tools-resin.sh.in
prefix=@prefix@

CCM_SERVLET_RESIN_CONF="${prefix}/share/ccm-tools/server/resin/bin/resin-conf-gen"

export CCM_SERVLET_RESIN_CONF
