# ccm-cms-types-formsectionitem/debian/prerm
#! /bin/sh
# prerm script for ccm-cms-types-formsectionitem
#
# see: dh_installdeb(1)
#
# Author: Berkan Eskikaya <berkan@runtime-collective.com>, 2004
#
# $Id: prerm,v 1.1.1.1 2004/11/12 09:30:06 fabrice Exp $
 
set -e

. /usr/share/debconf/confmodule
. /usr/share/ccm-tools/lib/shellmodule

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

db_get_ccm_settings ; db_stop

get_pkgname_from_arg "$0" 
ccm_package_name="$RET"
get_pkgversion_from_dpkg "$ccm_package_name"
ccm_package_version="$RET"

case "$1" in
    remove|purge|upgrade|deconfigure) 

	ccm_update_ccm_classpath remove $ccm_package_name $ccm_package_version
	ccm_update_ccm_webapps   remove $ccm_package_name $ccm_package_version
    ;;

    failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

if [ "$1" = "remove" -a "$ccm_package_name" = "ccm-core" -a -x "/etc/init.d/ccm" ]; then
  invoke-rc.d ccm stop || true
fi

 
# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.
 
#DEBHELPER#
 
exit 0

