# tools-legacy/tools/etc/profile.d/ccm-tools.sh.in
prefix=@prefix@

CCM_HOME=${prefix}/share/ccm
CCM_TOOLS_HOME=${prefix}/share/ccm-tools

export CCM_HOME
export CCM_TOOLS_HOME



