# tools-legacy/tools/etc/profile.d/Makefile.am
# Process with automake

AUTOMAKE_OPTIONS = foreign no-dependencies 1.4

profiledir = $(sysconfdir)/profile.d
profile_SCRIPTS = $(wildcard *.sh)

# End of files
